# deploy/soralink.service
[Unit]
Description=Soralink Relay Server
After=network.target

[Service]
Type=simple
User=soralink
ExecStart=/usr/local/bin/soralink-server --config /etc/soralink/server.yaml
Restart=always
RestartSec=5
LimitNOFILE=65536

# セキュリティ設定
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadOnlyPaths=/etc/soralink

[Install]
WantedBy=multi-user.target
